# backend/crontab
# Runs every minute; scripts/send_reminders.py decides what's actually due.
# Output goes to the container's stdout (PID 1) so `docker logs` shows it.
* * * * * root cd /app && /bin/uv run python scripts/send_reminders.py >> /proc/1/fd/1 2>&1
